# cloudflare/simple-ddns.timer.example
[Unit]
Description=Simple DDNS Timer

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
